# config/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

; graceful timeout should always be lower than stopwaitsecs to avoid orphan gunicorn workers.
[program:iskconDev-frappe-web]
command=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/env/bin/gunicorn -b 127.0.0.1:8000 -w 17 --max-requests 5000 --max-requests-jitter 500 -t 120 --graceful-timeout 30 frappe.app:application --preload
priority=4
autostart=true
autorestart=true
stdout_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/web.log
stderr_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/web.error.log
stopwaitsecs=40
killasgroup=true
user=yhpatt10
directory=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/sites
startretries=10

[program:iskconDev-frappe-schedule]
command=/Library/Frameworks/Python.framework/Versions/3.11/bin/bench schedule
priority=3
autostart=true
autorestart=true
stdout_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/schedule.log
stderr_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/schedule.error.log
user=yhpatt10
directory=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev
startretries=10



[program:iskconDev-frappe-short-worker]
command=/Library/Frameworks/Python.framework/Versions/3.11/bin/bench worker --queue short,default
priority=4
autostart=true
autorestart=true
stdout_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/worker.log
stderr_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/worker.error.log
user=yhpatt10
stopwaitsecs=360
directory=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev
killasgroup=true
numprocs=1
process_name=%(program_name)s-%(process_num)d
startretries=10

[program:iskconDev-frappe-long-worker]
command=/Library/Frameworks/Python.framework/Versions/3.11/bin/bench worker --queue long,default,short
priority=4
autostart=true
autorestart=true
stdout_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/worker.log
stderr_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/worker.error.log
user=yhpatt10
stopwaitsecs=1560
directory=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev
killasgroup=true
numprocs=1
process_name=%(program_name)s-%(process_num)d
startretries=10





[program:iskconDev-redis-cache]
command=/opt/homebrew/bin/redis-server /Users/yhpatt10/Documents/erpnextIskcon/iskconDev/config/redis_cache.conf
priority=1
autostart=true
autorestart=true
stdout_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/redis-cache.log
stderr_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/redis-cache.error.log
user=yhpatt10
directory=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/sites
startretries=10

[program:iskconDev-redis-queue]
command=/opt/homebrew/bin/redis-server /Users/yhpatt10/Documents/erpnextIskcon/iskconDev/config/redis_queue.conf
priority=1
autostart=true
autorestart=true
stdout_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/redis-queue.log
stderr_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/redis-queue.error.log
user=yhpatt10
directory=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/sites
startretries=10



[program:iskconDev-node-socketio]
command=/Users/yhpatt10/.nvm/versions/node/v22.3.0/bin/node /Users/yhpatt10/Documents/erpnextIskcon/iskconDev/apps/frappe/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/node-socketio.log
stderr_logfile=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev/logs/node-socketio.error.log
user=yhpatt10
directory=/Users/yhpatt10/Documents/erpnextIskcon/iskconDev
startretries=10


[group:iskconDev-web]
programs=iskconDev-frappe-web,iskconDev-node-socketio




[group:iskconDev-workers]
programs=iskconDev-frappe-schedule,iskconDev-frappe-short-worker,iskconDev-frappe-long-worker




[group:iskconDev-redis]
programs=iskconDev-redis-cache,iskconDev-redis-queue


[group:frappe]
